# deployment/systemd/bonbon-safety.service
[Unit]
Description=BonBon Safety Supervisor Runtime
After=docker.service network-online.target
Wants=docker.service network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
EnvironmentFile=/etc/bonbon/bonbon.env
WorkingDirectory=/opt/bonbon/current
ExecStart=/usr/bin/docker compose -f docker-compose.robot.yml up -d safety
ExecStop=/usr/bin/docker compose -f docker-compose.robot.yml stop safety
TimeoutStartSec=120

[Install]
WantedBy=multi-user.target
