# deployment/systemd/bonbon-speech.service
[Unit]
Description=BonBon Speech Runtime
After=bonbon-core.service

[Service]
Type=oneshot
RemainAfterExit=yes
EnvironmentFile=/etc/bonbon/bonbon.env
WorkingDirectory=/opt/bonbon/current
ExecStart=/usr/bin/docker compose -f docker-compose.robot.yml up -d speech
ExecStop=/usr/bin/docker compose -f docker-compose.robot.yml stop speech
TimeoutStartSec=120

[Install]
WantedBy=multi-user.target
