# deployment/systemd/bonbon-perception.service
[Unit]
Description=BonBon Perception Runtime
After=bonbon-core.service

[Service]
Type=oneshot
RemainAfterExit=yes
EnvironmentFile=/etc/bonbon/bonbon.env
WorkingDirectory=/opt/bonbon/current
ExecStart=/usr/bin/docker compose -f docker-compose.robot.yml up -d perception
ExecStop=/usr/bin/docker compose -f docker-compose.robot.yml stop perception
TimeoutStartSec=120

[Install]
WantedBy=multi-user.target
