# deployment/systemd/bonbon-core.service
[Unit]
Description=BonBon Core ROS2 Runtime
After=docker.service network-online.target
Wants=docker.service network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
EnvironmentFile=/etc/bonbon/bonbon.env
WorkingDirectory=/opt/bonbon/current
ExecStart=/usr/bin/docker compose -f docker-compose.robot.yml up -d core
ExecStop=/usr/bin/docker compose -f docker-compose.robot.yml stop core
TimeoutStartSec=120

[Install]
WantedBy=multi-user.target
